# sys/wszaudio.service
[Unit]
Description=wszaudio
 
[Service]
ExecStart=/home/wsz/.dot/sys/audio.sh
 
[Install]
WantedBy=multi-user.target
